# rootfs/etc/s6-overlay/s6-rc.d/openbox/run
#!/usr/bin/env sh

DISPLAY=:0 exec /usr/bin/openbox
