# rootfs/etc/s6-overlay/s6-rc.d/kcachegrind/run
#!/usr/bin/env sh

export $(/usr/bin/dbus-launch)
export DISPLAY=:0
export HOME=/data

exec /usr/bin/kcachegrind
